# docker/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.err.log
autorestart=true
